# packages/odoopbx/odoopbx-0.56.tar.gz/odoopbx-0.56/salt/roots/agent/agent.service
[Unit]
Description=OdooPBX Agent
After=network.target

[Service]
Environment=LANG=en_US.UTF-8
Type=simple
User=root
ExecStart=/usr/bin/env odoopbx run agent
Restart=on-failure
RestartSec=10
RemainAfterExit=no
StandardOutput=journal
StandardError=journal
SyslogIdentifier=odoopbx-agent

[Install]
WantedBy=multi-user.target
