# packages/mbsync-watcher/mbsync-watcher-0.1.3.tar.gz/mbsync-watcher-0.1.3/mbsync-watcher.service
[Unit]
Description=mbsync watcher

[Service]
Type=simple
ExecStart=%h/.local/bin/mbsync_watcher idle

[Install]
WantedBy=default.target
